# linux/profile.d/nvidia-conf-gnome.sh
# Nvidia tweaks for Gnome DE, allowing to run 2+ monitors with different refresh rates.
export CLUTTER_DEFAULT_FPS=144
