# linux/profile.d/nvidia-params.sh
# Nvidia tweaks for Gnome DE on XOrg, allowing to run 2+ monitors with different refresh rates.
export CLUTTER_DEFAULT_FPS=144
export __GL_SYNC_DISPLAY_DEVICE=DP-0
