# linux/etc/profile.d/90-ibus-lag-fix.sh
# Enable ibus hybrid asynchronous processing
export IBUS_ENABLE_SYNC_MODE=2
