# qt-match-gtk-scheme/.config/profile.d/10-qt-match-gtk-scheme.sh
# https://wiki.archlinux.org/title/Uniform_look_for_Qt_and_GTK_applications#QGtkStyle
export QT_QPA_PLATFORMTHEME=gtk3
