# mpd/.config/profile.d/10-set-mpd-address.sh
# Sets the mpd address so any client can automatically connect to it without
# manual configuration.
export MPD_HOST="$XDG_RUNTIME_DIR/mpd.socket"
