# shell-flatpak-path/.config/profile.d/10-flatpak-path.sh
export PATH="$PATH:/var/lib/flatpak/exports/bin"
