# system-nvidia-nvenc/etc/udev/rules.d/70-nvidia-nvenc.rules
ACTION=="add", DEVPATH=="/bus/pci/drivers/nvidia", RUN+="/usr/bin/nvidia-modprobe -c 0 -u"
